# crontab
# m h  dom mon dow   command
0 1 * * 1 /cronjobs/clean_database.sh
# end of file
